# bin/LiteSH.service
[Unit]
Description=Daemon_LiteSH


[Service]
Type=simple
User=root
Group=root
ExecStart=/home/wanderer/OS/CW_OS/bin/main
ExecReload=/home/wanderer/OS/CW_OS/bin/main
KillMode=control-group
TimeoutStartSec=120
OOMScoreAdjust=-1000


[Install]
WantedBy=multi-user.target
